# systemd/muslimtify.service
[Unit]
Description=Muslimtify prayer notification daemon
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/muslimtify daemon run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
